# .config/autostart/dunst.desktop
[Desktop Entry]
Exec=/usr/bin/dunst
Name=dunst
Type=Application
