# .config/autostart/kwalletd5.desktop
[Desktop Entry]
Exec=/bin/kwalletd5
Name=kwalletd5
Type=Application
